# templates/default/systemd/flannel.service.erb
[Unit]
Description=Flannel software-defined networking service
Documentation=https://coreos.com/flannel
After=network.target

[Service]
Type=notify
ExecStartPre=<%= @etcdctl_command %>
ExecStart=<%= @flanneld_command %>

[Install]
WantedBy=multi-user.target
